# rootfs/etc/services.d/cron/run
#!/usr/bin/with-contenv bashio

exec crond -f
